# data/export/systemd/process_master.target.erb
[Unit]
StopWhenUnneeded=false

[Install]
WantedBy=<%= app %>.target
WantedBy=multi-user.target